# 1_os/roles/archlinux_base/files/pacman-download-updates.service
[Unit]
Description=Download pacman updates
After=network-online.target

[Service]
Type=oneshot
ExecStart=pacman -Syu --downloadonly --noconfirm

